# scripts/runic-server.service
[Unit]
Description=Runic Firewall Control Plane
Documentation=https://github.com/ubenmackin/runic
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=runic
Group=runic
WorkingDirectory=/opt/runic
ExecStart=/opt/runic/dist/runic-server
Restart=on-failure
RestartSec=10s

# Secrets are loaded from the database at runtime

# Environment variables for non-secret configuration
Environment=RUNIC_DB_PATH=/opt/runic/data/runic.db
Environment=RUNIC_DOWNLOADS_DIR=/opt/runic/downloads
Environment=RUNIC_CERT_FILE=/opt/runic/cert/cert.pem
Environment=RUNIC_KEY_FILE=/opt/runic/cert/key.pem
Environment=RUNIC_PORT=60443
Environment=RUNIC_CONTROL_PLANE=localhost:60443

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/opt/runic/data /opt/runic/cert /opt/runic/downloads
ReadOnlyPaths=/etc/ssl/certs

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=runic-server

[Install]
WantedBy=multi-user.target
