# scripts/runic-agent.service
[Unit]
Description=Runic Firewall Agent
After=network-online.target
Wants=network-online.target
# Start after iptables/netfilter-persistent if present
After=netfilter-persistent.service
After=iptables.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/runic-agent --config /etc/runic-agent/config.json
# Restart=always ensures the agent restarts after a clean SIGTERM exit (e.g., during
# self-update via systemctl stop). To manually stop the agent, use:
# systemctl stop runic-agent && systemctl disable runic-agent
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=runic-agent

# Security hardening
NoNewPrivileges=false
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/etc/runic-agent /var/log

[Install]
WantedBy=multi-user.target
